# deploy/sherpa-staging.service
[Unit]
Description=SHERPA Staging
After=network.target

[Service]
Type=notify
User=dennis.kittrell
Group=dennis.kittrell
WorkingDirectory=/home/dennis.kittrell/sherpa-staging/code
EnvironmentFile=/home/dennis.kittrell/.percona-portal.env
EnvironmentFile=/home/dennis.kittrell/sherpa-staging/env.conf
ExecStart=/home/dennis.kittrell/sherpa-staging/venv/bin/gunicorn \
    --bind 0.0.0.0:3001 \
    --workers 1 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
