# deploy/sherpa-staging.service
[Unit]
Description=SHERPA Staging
After=network.target

[Service]
Type=simple
User=dennis.kittrell
Group=dennis.kittrell
WorkingDirectory=/home/dennis.kittrell/sherpa-staging/code
EnvironmentFile=/home/dennis.kittrell/.percona-portal.env
EnvironmentFile=/home/dennis.kittrell/sherpa-staging/env.conf
ExecStartPre=-/bin/rm -f /home/dennis.kittrell/sherpa-staging/code/gunicorn.ctl
ExecStart=/home/dennis.kittrell/sherpa-staging/venv/bin/gunicorn \
    --bind 0.0.0.0:3001 \
    --workers 1 \
    --timeout 120 \
    --no-control-socket \
    --access-logfile - \
    --error-logfile - \
    --graceful-timeout 5 \
    wsgi:app
ExecStopPost=-/bin/rm -f /home/dennis.kittrell/sherpa-staging/code/gunicorn.ctl
TimeoutStopSec=10
KillMode=mixed
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
